# supervisor.conf
[program:gunicorn]
command=gunicorn "app:create_app()" -w 4